# scripts/worldserver.service
[Unit]
Description=AzerothCore Worldserver
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=azerothcore
WorkingDirectory=/srv/azerothcore-wotlk-playerbots
ExecStart=/bin/screen -S worldserver -D -m /srv/azerothcore-wotlk-playerbots/acore.sh run-worldserver

[Install]
WantedBy=multi-user.target
